# test/e2e/quadlet/build.multiple-tags.container
## assert-podman-final-args localhost/imagename:v1

[Container]
Image=multiple-tags.build
