# test/e2e/quadlet/cgroups-mode.container
## assert-podman-args "--cgroups" "no-conmon"

[Container]
Image=localhost/imagename
CgroupsMode=no-conmon
