# test/e2e/quadlet/disableselinux.container
## assert-podman-args "--security-opt" "label:disable"

[Container]
Image=localhost/imagename
SecurityLabelDisable=true
