# test/e2e/quadlet/emptyline.container
[Container]
Exec=true \

# must have a blank line above, but this line can be anything (including another blank line)
