# test/e2e/quadlet/entrypoint-empty.container
## assert-podman-final-args localhost/imagename
## assert-podman-args "--entrypoint="

[Container]
Image=localhost/imagename
Entrypoint=
