# test/e2e/quadlet/entrypoint.container
## assert-podman-final-args localhost/imagename
## assert-podman-args "--entrypoint=top"

[Container]
Image=localhost/imagename
Entrypoint=top
