# test/e2e/quadlet/env-host.container
## assert-podman-final-args localhost/imagename
## assert-podman-args --env-host

[Container]
Image=localhost/imagename
EnvironmentHost=yes
