# test/e2e/quadlet/env-file-ignore.container
[Container]
Image=alpine
EnvironmentFile=-/tmp/missing-quadlet-file.env

## assert-podman-args "--env-file"
## assert-podman-args "-/tmp/missing-quadlet-file.env"
