# test/e2e/quadlet/env-host-false.container
## assert-podman-final-args localhost/imagename
## assert-podman-args --env-host=false

[Container]
Image=localhost/imagename
EnvironmentHost=no
