# test/e2e/quadlet/group.container
## assert-failed
## assert-stderr-contains "Group set without User"

[Container]
Image=localhost/imagename
Group=foobar
