# test/e2e/quadlet/hostname.container
[Container]
Image=localhost/imagename
## assert-podman-args "--hostname" "quadlet-host"
HostName=quadlet-host
