# test/e2e/quadlet/httpproxy-true.container
## assert-podman-final-args localhost/imagename
## assert-podman-args --http-proxy

[Container]
Image=localhost/imagename
HttpProxy=true
