# test/e2e/quadlet/httpproxy-false.container
## assert-podman-final-args localhost/imagename
## assert-podman-args --http-proxy=false

[Container]
Image=localhost/imagename
HttpProxy=false
