# test/e2e/quadlet/image.container
## assert-podman-final-args imagename

[Container]
Image=imagename
