# test/e2e/quadlet/image.container
## assert-podman-final-args localhost/imagename

[Container]
Image=localhost/imagename
