# test/e2e/quadlet/invalid-mount.container
## assert-failed
## assert-stderr-contains "invalid quadlet mount format"

[Container]
Image=localhost/imagename
Mount=/invalid:/mount:ro
