# test/e2e/quadlet/logdriver.container
## assert-podman-args "--log-driver" "journald"

[Container]
Image=localhost/imagename
LogDriver=journald
