# test/e2e/quadlet/login.quadlet.container
## assert-podman-args "--authfile" "/etc/certs/auth.json"

[Container]
Image=localhost/imagename
AuthFile=authfile.login
