# test/e2e/quadlet/memory.container
## assert-podman-final-args localhost/imagename
## assert-podman-args --memory 20g

[Container]
Image=localhost/imagename
Memory=20g
