# test/e2e/quadlet/mount-source-missing.container
## assert-failed
## assert-stderr-contains "source cannot be empty"

[Container]
Image=localhost/imagename
Mount=/home:/home
