# test/e2e/quadlet/name.container
## assert-podman-args "--name=foobar"

[Container]
Image=imagename
ContainerName=foobar
