# test/e2e/quadlet/name.container
## assert-podman-args "--name" "foobar"

[Container]
Image=localhost/imagename
ContainerName=foobar
