# test/e2e/quadlet/nestedselinux.container
## assert-podman-args "--security-opt" "label=nested"

[Container]
Image=localhost/imagename
SecurityLabelNested=true
