# test/e2e/quadlet/network.quadlet.container
## assert-podman-args "--network=systemd-basic"
## assert-key-is "Unit" "Requires" "systemd-basic-network.service"
## assert-key-is "Unit" "After" "systemd-basic-network.service"

[Container]
Image=localhost/imagename
Network=basic.network
