# test/e2e/quadlet/network.container
## assert-podman-args "--network=host"

[Container]
Image=localhost/imagename
Network=host
