# test/e2e/quadlet/noimage.container
## assert-failed
## assert-stderr-contains "No Image key specified"

[Container]
