# test/e2e/quadlet/noimage.container
## assert-failed
## assert-stderr-contains "no Image key specified"

[Container]
