# test/e2e/quadlet/noremapuser.container
## !assert-podman-args --uidmap
## !assert-podman-args --gidmap

[Container]
Image=imagename
RemapUsers=no
