# test/e2e/quadlet/notify-healthy.container
## assert-podman-args "--sdnotify=healthy"

[Container]
Image=localhost/imagename
Notify=healthy
