# test/e2e/quadlet/notify.container
## assert-podman-args "--sdnotify=container"

[Container]
Image=imagename
Notify=yes
