# test/e2e/quadlet/pod.not-found.container
## assert-failed
## assert-stderr-contains "quadlet pod unit not-found.pod does not exist"

[Container]
Image=localhost/imagename
Pod=not-found.pod
