# test/e2e/quadlet/pod.non-quadlet.container
## assert-failed
## assert-stderr-contains "pod test-pod is not Quadlet based"

[Container]
Image=localhost/imagename
Pod=test-pod
