# test/e2e/quadlet/readonly-tmpfs.container
## assert-podman-args "--read-only-tmpfs"
## assert-podman-args "--read-only"

[Container]
Image=localhost/imagename
ReadOnly=yes
ReadOnlyTmpfs=yes
