# test/e2e/quadlet/remap-keep-id.container
## assert-podman-args --userns keep-id

[Container]
Image=localhost/imagename
RemapUsers=keep-id
