# test/e2e/quadlet/remap-auto.container
## assert-podman-args --userns=auto

[Container]
Image=localhost/imagename
RemapUsers=auto
