# test/e2e/quadlet/rootfs.container
## assert-podman-final-args "--rootfs" "/var/lib/foobar"

[Container]
Rootfs=/var/lib/foobar
