# test/e2e/quadlet/root.container
## assert-podman-global-args "--root=/var/foobar/storage"

[Container]
Root=/var/foobar/storage
