# test/e2e/quadlet/seccomp.container
## assert-podman-args --security-opt seccomp=unconfined

[Container]
Image=localhost/imagename
SeccompProfile=unconfined
