# test/e2e/quadlet/service-name.container
## assert-podman-final-args localhost/imagename

[Container]
ServiceName=basic
Image=localhost/imagename
