# test/e2e/quadlet/shmsize.container
## assert-podman-args "--shm-size" "5g"

[Container]
Image=localhost/imagename
ShmSize=5g
