# test/e2e/quadlet/shortname.container
## assert-stderr-contains "not a fully qualified image name"

[Container]
Image=shortname
