# test/e2e/quadlet/stopsignal.container
## assert-podman-args "--stop-signal" "SIGKILL"

[Container]
Image=localhost/imagename
StopSignal=SIGKILL
