# test/e2e/quadlet/stoptimeout.container
## assert-podman-args "--stop-timeout" "20"

[Container]
Image=localhost/imagename
StopTimeout=20
