# test/e2e/quadlet/timezone.container
## assert-podman-args --tz=foo

[Container]
Image=imagename
Timezone=foo
