# test/e2e/quadlet/umask.container
## assert-podman-args "--umask" "0002"

[Container]
Image=localhost/imagename
Umask=0002
