# test/e2e/quadlet/unit-after-override.container
## assert-last-key-is-regex "Unit" "After" "^$"

[Unit]
After=

[Container]
Image=localhost/imagename
