# test/e2e/quadlet/user.container
## assert-podman-final-args localhost/imagename
## assert-podman-args "--user" "998:999"

[Container]
Image=localhost/imagename
User=998
Group=999
