# test/e2e/quadlet/user.container
## assert-podman-final-args imagename
## assert-podman-args "--user" "998:999"

[Container]
Image=imagename
User=998
Group=999
