# test/e2e/quadlet/userns.container
## assert-podman-args "--userns" "keep-id"

[Container]
Image=localhost/imagename
UserNS=keep-id
