# test/e2e/quadlet/user2.container
## assert-podman-final-args localhost/imagename
## assert-podman-args "--user=%U:%G"

[Container]
Image=localhost/imagename
User=%U
Group=%G
