# test/e2e/quadlet/volatiletmp-readonly.container
## !assert-podman-args "--tmpfs" "/tmp:rw,size=512M,mode=1777"

[Container]
Image=localhost/imagename
VolatileTmp=yes
ReadOnly=yes
