# test/e2e/quadlet/workingdir.container
## assert-podman-final-args localhost/imagename
## assert-podman-args "--workdir" "%h"

[Container]
Image=localhost/imagename
WorkingDir=%h
