# contrib/podmanimage/machine/etc/profile.d/docker-host.sh
export DOCKER_HOST="unix://$(podman info -f "{{.Host.RemoteSocket.Path}}")"
