# contrib/systemd/system/podman-kube@.service.in
[Unit]
Description=A template for running K8s workloads via podman-play-kube
Documentation=man:podman-play-kube(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=%t/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=never
TimeoutStopSec=70
ExecStart=@@PODMAN@@ play kube --replace --service-container=true %I
ExecStop=@@PODMAN@@ play kube --down %I
Type=notify
NotifyAccess=all

[Install]
WantedBy=default.target
