# contrib/systemd/system/podman-image-prune.timer
[Unit]
Description=Podman image prune timer

[Timer]
OnCalendar=daily
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
