# contrib/systemd/system/podman-auto-update.service.in
[Unit]
Description=Podman auto-update service
Documentation=man:podman-auto-update(1)
Wants=@@NETWORK_ONLINE_UNIT@@
After=@@NETWORK_ONLINE_UNIT@@

[Service]
Type=oneshot
ExecStart=@@PODMAN@@ auto-update
ExecStartPost=@@PODMAN@@ image prune -f

[Install]
WantedBy=default.target
