# contrib/varlink/io.podman.service
[Unit]
Description=Podman Remote API Service
Requires=io.podman.socket
After=io.podman.socket
Documentation=man:podman-varlink(1)

[Service]
Type=simple
ExecStart=/usr/bin/podman varlink unix:%t/podman/io.podman
KillMode=none

[Install]
WantedBy=multi-user.target
Also=io.podman.socket
